# tools/commissioner_thci/commissionerd.service.in
[Unit]
Description=OT-commissioner daemon.

[Service]
Type=simple
ExecStart=@PYTHON_BIN@ -u @COMMISSIONERD_PY_PATH@ -c @COMMISSIONER_CLI_PATH@

[Install]
WantedBy=multi-user.target
